# dot_config/autostart/joy_box_panel.desktop
[Desktop Entry]
Type=Application
Exec=/home/pi/joy_box_lcd_panel/joyBoxLcdPanTiltStart.sh

